# ops/systemd/sku-owl-autopatch.service
[Unit]
Description=SKU Owl Auto Patch and Rollback
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/sku-owl
ExecStart=/usr/bin/env bash /opt/sku-owl/ops/auto-patch-and-rollback.sh

